# debbuild/debian/prerm
#!/bin/sh

systemctl disable gowiki.service >/dev/null
systemctl stop gowiki.service >/dev/null