# debbuild/debian/prerm
#!/bin/sh

set -e

deb-systemd-helper disable gowiki.service >/dev/null
deb-systemd-helper stop gowiki.service >/dev/null